# kiosk shell/agent/systemd/kiosk-agent-queue.service
[Unit]
Description=Kiosk Agent Queue Runner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/kiosk-agent/agent/agent.py --config /etc/kiosk-agent/config.json --run-queue
WorkingDirectory=/opt/kiosk-agent
User=root
Environment=KIOSK_REPO_PATH=/opt/kiosk-app
StandardOutput=journal
StandardError=journal
