# license-check.service
[Unit]
Description=Kiosk License Check
Before=graphical.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/license-check.sh
User=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target