# 03_TASKLIST/security_licensing/systemd/alphaos-license-check.service
[Unit]
Description=AlphaOS License Validation
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/alphaos/license-url
Documentation=man:alphaos-license-check(8)

[Service]
Type=oneshot
User=alphaos-licd
Group=alphaos-licd
ExecStart=/opt/alphaos/license-check.sh
AmbientCapabilities=
CapabilityBoundingSet=
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallFilter=@system-service
ReadOnlyPaths=/
ReadWritePaths=/var/lib/alphaos /run/alphaos
StateDirectory=alphaos
Restart=no

[Install]
WantedBy=multi-user.target
